# deploy/ansible/httpd/template/httpd.sysv.j2
#!/bin/bash

# chkconfig: 2345 64 36
# description: httpd server script
HTTPD=/usr/local/httpd/bin/httpd
PID_FILE=/usr/local/httpd/logs/httpd.pid

case $1 in
"start")
    if test -e $PID_FILE
    then
        echo 'httpd is started ...'
    else
        echo 'httpd will start ...'
        $HTTPD
    fi
;;

"stop")
    if test -e $PID_FILE
    then
        PID=`cat $PID_FILE`
        kill $PID
    else
        echo 'httpd is stoped'
    fi
;;
"status")
    if test -e $PID_FILE
    then
        echo 'httpd has been started '
    else
        echo 'httpd is stoped'
    fi
;;
*)
    echo "not support option $1"
;;
esac